# sv/3g-route/finish
#!/bin/bash
sleep 1
exec 2>&1
/sbin/route del -net default metric 0 dev ppp0
sysctl net.ipv4.ip_forward=0
iptables -t nat -X
iptables -t nat -F
